# install/linux/roninKB-daemon.service
[Unit]
Description=RoninKB Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hhkb-daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
